# rootfs/etc/services.d/kres-cache-gc/run
#!/command/execlineb -P

/usr/sbin/kres-cache-gc -c /var/cache/knot-resolver -d 1000
