# rootfs/etc/services.d/kresd/run
#!/command/execlineb -P

cd /var/lib/knot-resolver
/usr/sbin/kresd -c /etc/knot-resolver/kresd.conf -n
